# etc/profile.d/arduino.sh
export ARDUINO_HOME="/opt/arduino/default"

export PATH="$ARDUINO_HOME:$PATH"
