# etc/profile.d/eclipse.sh
export ECLIPSE_HOME="/opt/eclipse/default"

export PATH="$ECLIPSE_HOME:$PATH"
